# etc/services.d/nginx/run
#!/usr/bin/execlineb -P
exec /usr/sbin/nginx -g "daemon off;"